# make-snapshots.target
[Unit]
Description=Make Snapshots
Requires=multi-user.target
